# azkaban-webserver/root/etc/services.d/azkaban-webserver/run
#!/usr/bin/with-contenv /bin/sh

runuser -l %AZKABAN_USER% -c 'cd /opt/azkaban-web-server && bin/azkaban-web-start.sh'