# azkaban-execserver/root/etc/services.d/azkaban-execserver/run
#!/usr/bin/with-contenv /bin/sh

runuser -l %AZKABAN_USER% -c 'cd /opt/azkaban-exec-server && bin/azkaban-executor-start.sh'